# dev-container/scripts/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:code-server]
command=/usr/bin/code-server --bind-addr 0.0.0.0:8080 --auth none /workspace
user=devuser
environment=HOME="/home/devuser",XDG_CONFIG_HOME="/home/devuser/.config",XDG_DATA_HOME="/home/devuser/.local/share"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:filebrowser]
command=/usr/local/bin/filebrowser -a 0.0.0.0 -p 8443 -r /workspace --baseurl /files --noauth -d /home/devuser/.filebrowser.db
user=devuser
environment=HOME="/home/devuser"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:gotty]
command=/usr/local/bin/gotty -w --permit-arguments --reconnect -a 0.0.0.0 -p 7681 /usr/local/bin/launch-cli.sh
user=devuser
environment=HOME="/home/devuser",TERM="xterm-256color"
directory=/workspace
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:t3]
command=/bin/bash -lc "cd /workspace && exec t3 serve --host 0.0.0.0 --port 3773 --no-browser"
user=devuser
environment=HOME="/home/devuser",TERM="xterm-256color"
directory=/workspace
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:t3-link]
command=/usr/bin/python3 /usr/local/bin/t3-link-server.py
user=devuser
environment=HOME="/home/devuser",TERM="xterm-256color"
directory=/workspace
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
