# dev-container/scripts/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1280x1024x24 -ac
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:fluxbox]
command=/usr/bin/fluxbox
user=devuser
environment=DISPLAY=":99",HOME="/home/devuser"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:idesk]
command=/usr/bin/idesk
user=devuser
environment=DISPLAY=":99",HOME="/home/devuser"
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:wallpaper]
command=/bin/bash -c "sleep 2 && DISPLAY=:99 feh --bg-fill /home/devuser/.fluxbox/wallpaper.png"
user=devuser
environment=DISPLAY=":99",HOME="/home/devuser"
autostart=true
autorestart=false
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -forever -nopw -shared -rfbport 5900
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:websockify]
command=/usr/bin/websockify --web /usr/share/novnc 6080 localhost:5900
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:code-server]
command=/usr/bin/code-server --bind-addr 0.0.0.0:8080 --auth none /workspace
user=devuser
environment=DISPLAY=":99",HOME="/home/devuser",XDG_CONFIG_HOME="/home/devuser/.config",XDG_DATA_HOME="/home/devuser/.local/share"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:filebrowser]
command=/usr/local/bin/filebrowser -a 0.0.0.0 -p 8443 -r /workspace --baseurl /files --noauth -d /home/devuser/.filebrowser.db
user=devuser
environment=HOME="/home/devuser"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:gotty]
command=/usr/local/bin/gotty -w --permit-arguments --reconnect -a 0.0.0.0 -p 7681 /usr/local/bin/launch-cli.sh
user=devuser
environment=HOME="/home/devuser",TERM="xterm-256color"
directory=/workspace
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
